# scripts/beacons.service
[Unit]
Description=Update beacons every 12 hours

[Service]
Type=simple
ExecStart=/var/lib/KST/beacons.sh

[Install]
WantedBy=multi-user.target

